# daemon/systemd/teranode-explorer-watchdog.service
[Unit]
Description=Teranode explorer listener watchdog (detect silent gossip-flow death and auto-restart)
After=network-online.target teranode-explorer-listener.service
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=/etc/bhs-watchdog/env
ExecStart=/usr/local/bin/teranode-explorer-watchdog.sh
